# winspotlightkde.timer
[Unit]
Description=WinSpotlightKDE

[Timer]
Unit=winspotlightkde.service
OnCalendar=*-*-* 09:05:00
Persistent=true

[Install]
WantedBy=timers.target
